# nomad.conf
description "Nomad"
author "Francis Tseng <@frnsys>"

start on runlevel [234]
stop on runlevel [0156]

chdir /some/dir
exec /some/dir/script.py
respawn
